# tomcat/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=/bin/bash /tomcat-start.sh
