# krb5-server/supervisord.conf
[supervisord]
nodaemon=true

[program:krb5kdc]
command=/usr/sbin/krb5kdc -n

[program:kadmind]
command=/usr/sbin/kadmind -nofork
