# spree-nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/bin/bash /nginx-start.sh
