# haproxy/supervisord.conf
[supervisord]
nodaemon=true

[program:syslog]
command=/usr/sbin/rsyslogd -n

[program:haproxy]
command=/bin/bash /docker-haproxy.sh
