# pgpool/supervisord.conf
[supervisord]
nodaemon=true

[program:pgpool]
command=/bin/bash /docker-pgpool.sh
